# systemd/keybase.buildplease.timer
[Unit]
Description=run keybase.buildplease.service every weekday at noon

[Timer]
OnCalendar=Mon-Fri 12:00

[Install]
WantedBy=timers.target
